# var/lib/dpkg/info/glances.conffiles
/etc/default/glances
/etc/glances/glances.conf
/etc/init.d/glances
