# var/lib/dpkg/info/glances.list
/.
/usr
/usr/lib
/usr/lib/python3
/usr/lib/python3/dist-packages
/usr/lib/python3/dist-packages/glances
/usr/lib/python3/dist-packages/glances/__main__.py
/usr/lib/python3/dist-packages/glances/__init__.py
/usr/lib/python3/dist-packages/glances/outputs
/usr/lib/python3/dist-packages/glances/outputs/__init__.py
/usr/lib/python3/dist-packages/glances/outputs/bottle
/usr/lib/python3/dist-packages/glances/outputs/bottle/base.tpl
/usr/lib/python3/dist-packages/glances/outputs/bottle/plugin_text.tpl
/usr/lib/python3/dist-packages/glances/outputs/bottle/plugin_table.tpl
/usr/lib/python3/dist-packages/glances/outputs/glances_bottle.py
/usr/lib/python3/dist-packages/glances/outputs/glances_colorconsole.py
/usr/lib/python3/dist-packages/glances/outputs/glances_bars.py
/usr/lib/python3/dist-packages/glances/outputs/static
/usr/lib/python3/dist-packages/glances/outputs/static/css
/usr/lib/python3/dist-packages/glances/outputs/static/css/style.css
/usr/lib/python3/dist-packages/glances/outputs/static/css/bootstrap.min.css
/usr/lib/python3/dist-packages/glances/outputs/static/css/normalize.css
/usr/lib/python3/dist-packages/glances/outputs/static/js
/usr/lib/python3/dist-packages/glances/outputs/static/js/modernizr.custom.js
/usr/lib/python3/dist-packages/glances/outputs/static/favicon.ico
/usr/lib/python3/dist-packages/glances/outputs/glances_curses.py
/usr/lib/python3/dist-packages/glances/exports
/usr/lib/python3/dist-packages/glances/exports/glances_export.py
/usr/lib/python3/dist-packages/glances/exports/__init__.py
/usr/lib/python3/dist-packages/glances/exports/glances_history.py
/usr/lib/python3/dist-packages/glances/exports/glances_influxdb.py
/usr/lib/python3/dist-packages/glances/exports/glances_csv.py
/usr/lib/python3/dist-packages/glances/exports/glances_statsd.py
/usr/lib/python3/dist-packages/glances/plugins
/usr/lib/python3/dist-packages/glances/plugins/glances_memswap.py
/usr/lib/python3/dist-packages/glances/plugins/glances_system.py
/usr/lib/python3/dist-packages/glances/plugins/glances_psutilversion.py
/usr/lib/python3/dist-packages/glances/plugins/glances_diskio.py
/usr/lib/python3/dist-packages/glances/plugins/glances_cpu.py
/usr/lib/python3/dist-packages/glances/plugins/glances_processlist.py
/usr/lib/python3/dist-packages/glances/plugins/__init__.py
/usr/lib/python3/dist-packages/glances/plugins/glances_monitor.py
/usr/lib/python3/dist-packages/glances/plugins/glances_mem.py
/usr/lib/python3/dist-packages/glances/plugins/glances_raid.py
/usr/lib/python3/dist-packages/glances/plugins/glances_fs.py
/usr/lib/python3/dist-packages/glances/plugins/glances_percpu.py
/usr/lib/python3/dist-packages/glances/plugins/glances_hddtemp.py
/usr/lib/python3/dist-packages/glances/plugins/glances_network.py
/usr/lib/python3/dist-packages/glances/plugins/glances_processcount.py
/usr/lib/python3/dist-packages/glances/plugins/glances_now.py
/usr/lib/python3/dist-packages/glances/plugins/glances_alert.py
/usr/lib/python3/dist-packages/glances/plugins/glances_load.py
/usr/lib/python3/dist-packages/glances/plugins/glances_uptime.py
/usr/lib/python3/dist-packages/glances/plugins/glances_docker.py
/usr/lib/python3/dist-packages/glances/plugins/glances_batpercent.py
/usr/lib/python3/dist-packages/glances/plugins/glances_core.py
/usr/lib/python3/dist-packages/glances/plugins/glances_plugin.py
/usr/lib/python3/dist-packages/glances/plugins/glances_help.py
/usr/lib/python3/dist-packages/glances/plugins/glances_sensors.py
/usr/lib/python3/dist-packages/glances/core
/usr/lib/python3/dist-packages/glances/core/glances_logs.py
/usr/lib/python3/dist-packages/glances/core/glances_client_browser.py
/usr/lib/python3/dist-packages/glances/core/glances_actions.py
/usr/lib/python3/dist-packages/glances/core/glances_timer.py
/usr/lib/python3/dist-packages/glances/core/glances_processes.py
/usr/lib/python3/dist-packages/glances/core/__init__.py
/usr/lib/python3/dist-packages/glances/core/glances_snmp.py
/usr/lib/python3/dist-packages/glances/core/glances_logging.py
/usr/lib/python3/dist-packages/glances/core/glances_client.py
/usr/lib/python3/dist-packages/glances/core/glances_server.py
/usr/lib/python3/dist-packages/glances/core/glances_monitor_list.py
/usr/lib/python3/dist-packages/glances/core/glances_standalone.py
/usr/lib/python3/dist-packages/glances/core/glances_staticlist.py
/usr/lib/python3/dist-packages/glances/core/glances_config.py
/usr/lib/python3/dist-packages/glances/core/glances_webserver.py
/usr/lib/python3/dist-packages/glances/core/glances_main.py
/usr/lib/python3/dist-packages/glances/core/glances_autodiscover.py
/usr/lib/python3/dist-packages/glances/core/glances_stats.py
/usr/lib/python3/dist-packages/glances/core/glances_globals.py
/usr/lib/python3/dist-packages/glances/core/glances_password.py
/usr/lib/python3/dist-packages/Glances-2.3.egg-info
/usr/lib/python3/dist-packages/Glances-2.3.egg-info/dependency_links.txt
/usr/lib/python3/dist-packages/Glances-2.3.egg-info/entry_points.txt
/usr/lib/python3/dist-packages/Glances-2.3.egg-info/requires.txt
/usr/lib/python3/dist-packages/Glances-2.3.egg-info/top_level.txt
/usr/lib/python3/dist-packages/Glances-2.3.egg-info/PKG-INFO
/usr/share
/usr/share/doc
/usr/share/doc/glances
/usr/share/doc/glances/copyright
/usr/share/doc/glances/README.rst.gz
/usr/share/doc/glances/glances-doc.html
/usr/share/doc/glances/changelog.Debian.gz
/usr/share/doc/glances/images
/usr/share/doc/glances/images/per-cpu.png
/usr/share/doc/glances/images/logs.png
/usr/share/doc/glances/images/network.png
/usr/share/doc/glances/images/disconnected.png
/usr/share/doc/glances/images/fs.png
/usr/share/doc/glances/images/monitored.png
/usr/share/doc/glances/images/mem-wide.png
/usr/share/doc/glances/images/docker.png
/usr/share/doc/glances/images/diskio.png
/usr/share/doc/glances/images/header.png
/usr/share/doc/glances/images/processlist.png
/usr/share/doc/glances/images/sensors.png
/usr/share/doc/glances/images/screenshot-web.png
/usr/share/doc/glances/images/mem.png
/usr/share/doc/glances/images/cpu-wide.png
/usr/share/doc/glances/images/connected.png
/usr/share/doc/glances/images/screenshot-wide.png
/usr/share/doc/glances/images/screenshot.png
/usr/share/doc/glances/images/processlist-wide.png
/usr/share/doc/glances/images/glances-responsive-webdesign.png
/usr/share/doc/glances/images/processlist-top.png
/usr/share/doc/glances/images/screenshot-web2.png
/usr/share/doc/glances/images/raid.png
/usr/share/doc/glances/images/load.png
/usr/share/doc/glances/images/cpu.png
/usr/share/doc/glances/examples
/usr/share/doc/glances/examples/glances.conf
/usr/share/doc/glances/NEWS.Debian.gz
/usr/share/man
/usr/share/man/man1
/usr/share/man/man1/glances.1.gz
/usr/share/locale
/usr/share/locale/de
/usr/share/locale/de/LC_MESSAGES
/usr/share/locale/de/LC_MESSAGES/glances.mo
/usr/share/doc-base
/usr/share/doc-base/glances-manual
/usr/bin
/usr/bin/glances
/etc
/etc/init.d
/etc/init.d/glances
/etc/glances
/etc/glances/glances.conf
/etc/default
/etc/default/glances
