# var/lib/dpkg/info/ieee-data.list
/.
/var
/var/lib
/var/lib/ieee-data
/usr
/usr/share
/usr/share/ieee-data
/usr/share/ieee-data/.lastupdate
/usr/share/ieee-data/iab.txt
/usr/share/ieee-data/oui.txt
/usr/share/doc
/usr/share/doc/ieee-data
/usr/share/doc/ieee-data/copyright
/usr/share/doc/ieee-data/changelog.gz
/usr/share/doc/ieee-data/README.Debian
/usr/share/man
/usr/share/man/man8
/usr/share/man/man8/update-oui.8.gz
/usr/bin
/usr/bin/update-oui
/etc
/etc/cron.monthly
/etc/cron.monthly/ieee-data
/var/lib/ieee-data/.lastupdate
/var/lib/ieee-data/iab.txt
/var/lib/ieee-data/oui.txt
/var/lib/ieee-data/oui36.txt
/var/lib/ieee-data/mam.txt
